# deploy/systemd/sovereign-sensor-watchdog.service
[Unit]
Description=Sovereign Sensor Agent Freshness Watchdog
After=sovereign-sensor-ingest.service

[Service]
Type=oneshot
User=ssa
WorkingDirectory=/opt/sovereign-sensor-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sovereign-sensor-agent/.venv/bin/python /opt/sovereign-sensor-agent/scripts/watchdog.py
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
