# deploy/systemd/sovereign-sensor-ingest.service
[Unit]
Description=Sovereign Sensor Agent Serial Ingest
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dhuzard
WorkingDirectory=/home/dhuzard/sovereign-sensor-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/dhuzard/sovereign-sensor-agent/.venv/bin/python /home/dhuzard/sovereign-sensor-agent/scripts/read_serial.py --device /dev/ttyACM0 --baud 9600 --sensor-id arduino-ttyACM0
Restart=always
RestartSec=3
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
# Security: ingest needs serial device access but nothing else
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target