# deploy/systemd/sovereign-sensor-ingest.service
[Unit]
Description=Sovereign Sensor Agent Serial Ingest
After=network.target

[Service]
Type=simple
User=dhuzard
WorkingDirectory=/home/dhuzard/sovereign-sensor-agent
ExecStart=/home/dhuzard/sovereign-sensor-agent/.venv/bin/python /home/dhuzard/sovereign-sensor-agent/scripts/read_serial.py --device /dev/ttyACM0 --baud 9600 --sensor-id arduino-ttyACM0
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target