# deploy/systemd/sovereign-sensor-api.service
[Unit]
Description=Sovereign Sensor Agent API
Wants=network-online.target
After=network-online.target sovereign-sensor-ingest.service
Requires=sovereign-sensor-ingest.service

[Service]
Type=simple
User=dhuzard
WorkingDirectory=/home/dhuzard/sovereign-sensor-agent
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/home/dhuzard/sovereign-sensor-agent/.venv/bin/python /home/dhuzard/sovereign-sensor-agent/scripts/check_env.py --service api
ExecStart=/home/dhuzard/sovereign-sensor-agent/.venv/bin/python /home/dhuzard/sovereign-sensor-agent/scripts/api_server.py --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
# Security: API serves read-only HTTP; no device or serial access needed
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target