# deploy/systemd/sovereign-sensor-nanobot.service
[Unit]
Description=Sovereign Sensor Agent WhatsApp Stack
Wants=network-online.target
After=network-online.target sovereign-sensor-api.service
Requires=sovereign-sensor-api.service

[Service]
Type=simple
User=dhuzard
WorkingDirectory=/home/dhuzard/sovereign-sensor-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/dhuzard/sovereign-sensor-agent/deploy/nanobot/start_nanobot.sh up
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target